# packages/os/healthdog.service
[Unit]
Description=Send a Healthdog Ping

[Service]
Type=oneshot
RemainAfterExit=false
ExecStart=/bin/healthdog send-health-ping

[Install]
WantedBy=multi-user.target
