# packages/os/healthdog.timer
[Unit]
Description=Scheduled Healthdog Pings

[Timer]
# Don't run missed executions
Persistent=false
# Run 120 seconds after boot for the first time
OnBootSec=120
# Run every 6 hours thereafter #21600
OnUnitActiveSec=6
# File describing job to execute
Unit=healthdog.service

[Install]
WantedBy=timers.target
